# systemd/user/cloud-backup/cloud-backup.timer
[Unit]
Description=Run `cloud-backup` hourly after login

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
