# systemd/user/cloud-backup/cloud-backup.service
[Unit]
Description=rclone copy backup <-> cloud

[Service]
ExecStart=%h/dotfiles/systemd/user/cloud-backup/cloud-backup.sh
Type=oneshot
