# systemd/check-updates/check-updates.timer
[Unit]
Description=Run 'check-updates' hourly

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
