# systemd/check-updates/check-updates.service
[Unit]
Description=check updates

[Service]
Type=oneshot
ExecStart=pacman -Sy
